# crontab
* * * * * /usr/local/bin/python /app/main.py >> /var/log/auto_ria/test.log 2>&1
